# supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=uvicorn api.app:app --host 0.0.0.0 --port 8000

[program:bot]
command=python telegram_bot/bot.py